# power-outage.service
[Unit]
Description=LOE Power Outage Scraper → Home Assistant
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/home/YOUR_USERNAME/power_outages
ExecStart=/usr/bin/python3 /home/YOUR_USERNAME/power_outages/power_outage_scraper.py
Restart=on-failure
RestartSec=60
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
